# build/4.0.0/pkg/hsdpa/hsdpa.pkg.d/hsdpa/53-usb_wimax_huddle.rules
#####################################
# Huddle device block
#####################################

SUBSYSTEM=="usb", ACTION=="add", ATTRS{idVendor}=="19f2", ATTRS{idProduct}=="1700", RUN+="/bin/sh /usr/local/apache/qb/setuid/checkHuddlestate.sh up"

